# dist/medisana-bs.service
[Unit]
Description=Run Medisana bluetooth scale monitor
After=network-online.target hciuart.service bluetooth.target
Wants=network-online.target systemd-networkd-wait-online.service hciuart.service bluetooth.target

[Service]
WorkingDirectory=/opt/medisana-bs
ExecStartPre=//usr/bin/btmgmt le on
ExecStart=/opt/medisana-bs/medisana-bs.arm6
TimeoutStopSec=1s
StartLimitInterval=20
StartLimitBurst=5
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
